# ToDo/ToDo/deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for SITENAME"

start on net-device-up
stop on shutdown

respawn

setuid jaccorot
chdir /home/jaccorot/sites/SITENAME/source

exec gunicorn --bind unix:/tmp/SITENAME.socket todo.wsgi:application
